# yiya/crontab
0 21  *    * 1,4  /home/pi/git/home-cluster/yiya/start-backup.sh >> /home/pi/cronlog.log
