# donna/crontab
0    21  * * 0,3    /home/pi/git/home-cluster/donna/start-backup.sh
0    */8 * * *      sudo mount -a
*/30 *   * * *      /home/pi/git/home-cluster/donna/protonvpn-ping.sh
